# modules/puppet/templates/git-puppet.timer
[Unit]
Description=Schedule periodic Puppet runs

[Timer]
OnActiveSec=30
OnUnitActiveSec=<%= @interval %>

[Install]
WantedBy=multi-user.target
